# configs/supervisor-analytics.conf
[program:analytics]
command=/bin/bash /scripts/analytics.sh
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/null
stderr_logfile=/dev/null
priority=10
