# configs/supervisor-dashboard.conf
[program:dashboard]
command=node /opt/dashboard-server.js
directory=/opt
user=ubuntu
environment=HOME="/home/ubuntu",NODE_ENV="production"
priority=100
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/tmp/dashboard_stdout.log
stderr_logfile=/tmp/dashboard_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
