# configs/supervisor-claudecodeui.conf
[program:claudecodeui]
command=/usr/local/bin/claude-code-ui --port 10007
directory=/home/ubuntu/projects
user=ubuntu
environment=HOME="/home/ubuntu",NODE_ENV="production",PORT="10007"
priority=200
autostart=%(ENV_ENABLE_CLAUDECODEUI)s
autorestart=true
startretries=5
startsecs=15
stdout_logfile=/tmp/claudecodeui_stdout.log
stderr_logfile=/tmp/claudecodeui_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
